# deploy/imalink-core.service
[Unit]
Description=ImaLink Core API Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/imalink-core
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/.local/bin/uv run uvicorn service.main:app --host 127.0.0.1 --port 8765
Restart=always
RestartSec=10

# VIKTIG: Oppdater følgende før du kjører systemctl:
# 1. User=dinbruker (erstatt www-data)
# 2. Group=dinbruker (erstatt www-data)
# 3. WorkingDirectory=/home/dinbruker/imalink-core (din faktiske sti)
# 4. Environment="PATH=..." (oppdater hvis uv er installert et annet sted)

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
